# rc.local
#!/bin/sh -e

/home/api/forever-start.sh

exit 0
