# rc.local
#!/bin/sh -e

#ALIAS
source /etc/profile.d/00-alias.sh

#APIS
/home/api/forever-start.sh

exit 0
