# profile.d/00-alias.sh
alias ll='ls -alF'
alias targz='tar -zxvf'
